# container@.service
[Unit]
Description=Anchorage container service: %i
Documentation=man:anchorage(7)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=anchorage
Group=anchorage
WorkingDirectory=/var/lib/anchorage/%i
Environment=SERVICE_NAME=%i
Delegate=yes
TimeoutStopSec=120

ExecStart=/usr/lib/anchorage/run.sh
ExecStop=/usr/bin/podman-compose down
ExecReload=/usr/bin/podman-compose pull
ExecReload=/usr/bin/podman-compose up -d --force-recreate

[Install]
WantedBy=multi-user.target
